# roles/shipping/templates/shipping.service
[Unit]
Description=Shipping Service

[Service]
User=roboshop
Environment=CART_ENDPOINT=CART-{{ env }}.vishnureddy.online:8080
Environment=DB_HOST=MYSQL{{ env }}.vishnureddy.online
ExecStart=/bin/java -jar /app/shipping.jar
SyslogIdentifier=shipping

[Install]
WantedBy=multi-user.target
